# scripts/com.passive-agent.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.passive-agent.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/.venv/bin/passive-agent</string>
        <string>daily</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/data/reports/daily_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/data/reports/daily_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:__PROJECT_DIR__/.venv/bin</string>
    </dict>
</dict>
</plist>
